# service/app1/finish
#!/usr/bin/env sh

echo "app1 挂掉啦！..."
echo "app1 restart ok!"
