# service/app2/run
#!/usr/bin/env sh

echo "Started APP2 service..."

for i in {1..3}; do
    echo "这里是第2个程序！"
    sleep 1
done

echo "Oh no!我嗝屁了..." >&2 #3秒就死掉
exit 1
